# deployment/supervisor.conf
; ── Supervisor Configuration — Laravel Queue Workers ──────────
; File: /etc/supervisor/conf.d/lnc-queue.conf
; Install: sudo apt-get install supervisor

[program:lnc-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/lnc/backend/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/lnc-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600

[program:lnc-scheduler]
; This runs every minute — acts as the Laravel cron alternative via supervisor
process_name=%(program_name)s
command=bash -c "while true; do php /var/www/lnc/backend/artisan schedule:run >> /var/log/supervisor/lnc-scheduler.log 2>&1; sleep 60; done"
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/lnc-scheduler.log
stdout_logfile_maxbytes=5MB
